# cron.d/root
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
*/10 * * * * sh /etc/update.sh 
